# config/includes.chroot_after_packages/etc/xdg/autostart/xfce-fix.desktop
[Desktop Entry]
Type=Application
Name=XFCE Fix
Exec=/usr/local/bin/xfce-fix.sh
OnlyShowIn=XFCE;